# systemd/wg-panel.service
[Unit]
Description=WireGuard Server Panel
After=network.target

[Service]
User=root
WorkingDirectory=/opt/wg-panel
ExecStart=/opt/wg-panel/venv/bin/python -m gunicorn.app.wsgiapp --workers 2 --threads 4 --bind 0.0.0.0:5000 --access-logfile - --error-logfile - app:app
Restart=always
Environment=PYTHONUNBUFFERED=1
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
